# temp/reference/systemd_config/oak-detector.service
[Unit]
Description=OAK Detection System (GUI + Conda)
After=graphical.target network-online.target simple_can_setup.service
Wants=graphical.target network-online.target simple_can_setup.service
# 确保在系统完全启动后再启动服务
After=multi-user.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/oak2/temp
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/root/.Xauthority
ExecStartPre=/bin/sleep 10

# 方式1（推荐）：直接用 Conda 环境的 python 绝对路径
# 把 /root/miniconda3/envs/oak/bin/python 替换为你的实际 env/python
ExecStart=/root/miniconda3/envs/oak/bin/python /root/oak2/temp/dual_detectionv2.py

# 方式2（可选）：使用 conda run（保留注释，二选一）
# ExecStart=/root/miniconda3/bin/conda run --no-capture-output -n oak python /root/oak2/temp/dual_detectionv2.py

Restart=no
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target